# docker/supervisord.conf
; Mantis — supervisor config.
; Runs the Rust backend and a tiny static UI server inside one container.
; Both processes log to stdout so `docker logs` shows everything inline.

[supervisord]
nodaemon=true
user=mantis
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:mantis-backend]
command=/app/mantis-backend
directory=/app
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:mantis-ui]
command=python3 -m http.server 5173 --bind 0.0.0.0
directory=/app/ui
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
